# build/rootfs/etc/supervisor/conf.d/nginx.conf
[program:nginx]
priority=10
command=nginx -c /etc/nginx/nginx.conf -g 'daemon off;'
